# etc/profile.d/golang.sh
export GOPATH=$HOME/code/go
export GOROOT=$HOME/.asdf/installs/golang/1.9/go/
