# etc/profile.d/golang.sh
export GOPATH=$HOME/code/go
export GOROOT=~/.asdf/installs/golang/1.8.3/go/
