# etc/profile.d/postgres.sh
#!/bin/bash

export PG_HOST=localhost
export PG_PORT=5432
export PG_PASSWORD=
export PG_USER=
