# etc/profile.d/postgres.sh
#!/bin/bash

export PGHOST=localhost
export PGPORT=5432
export PGPASSWORD=
export PGUSER=$USER
export DPGUSER=
export DPGPASS=
